# ops/deploy/systemd/aionbd.service
[Unit]
Description=AIONBD Vector Database
Documentation=https://github.com/<org>/aionbd
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aionbd
Group=aionbd
EnvironmentFile=-/etc/aionbd/aionbd.env
ExecStart=/usr/local/bin/aionbd-server
Restart=on-failure
RestartSec=2
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
ReadWritePaths=/var/lib/aionbd

[Install]
WantedBy=multi-user.target
